# files/foo.service
[Unit]
Description=This always fails

[Install]
WantedBy=multi-user.target
